# etc/crontab
*/30 * * * * sleep $(expr $RANDOM \% 300) ; cd /SackConfig && ./generate.sh
